# supervisord.conf
[program:runserver]
command=/home/ubuntu/mazio/mazioenv/bin/python /home/ubuntu/mazio/tastygis/mozio/manage.py runserver 0:8000
stopsignal=KILL
killasgroup=true
